# srv/var/spool/cron/crontabs/root
*/5 *	* * *	/usr/bin/python3 /var/www/news-bot/src/MyTeleBot/newsbot_rss.py  >> /var/log/newsbot.rss.log 2>&1
